# mod_zynthian/systemd/performance.service
[Unit]
Description=Tune CPU perfoamnce

[Service]
Type=oneshot
ExecStart=/home/pi/bin/cpu-performance.sh

[Install]
WantedBy=multi-user.target
