# mod_pisound/systemd/pisound-btn.service
[Unit]
Description=pisound-btn server

[Service]
ExecStart=/usr/local/bin/pisound-btn

[Install]
WantedBy=multi-user.target
