# tmp/udev/50-mod-usbmidi.rules
ACTION=="add", KERNEL=="midiC*", RUN+="/bin/systemctl restart jack2"
ACTION=="remove", KERNEL=="midiC*", RUN+="/bin/systemctl restart jack2"
