# mod_zynthian/udev/50-mod-usbmidi.rules
ACTION=="add", KERNEL=="midiC*", RUN+="/home/pi/mod_midi_autoconnect.py", RUN+="/bin/systemctl restart mod-host"
